# systemd/nvidia_tweaks.service
[Unit]
Description=Nvidia Tweaks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nvidia_tweaks.sh

[Install]
WantedBy=multi-user.target
